# church-api/deploy/supervisor/church-api.conf
; ---------------------------------------------------------------------------
; Supervisor — church-api realtime + async processes (CHR-158)
; ---------------------------------------------------------------------------
; Keeps the Reverb WebSocket server, the queue workers and the scheduler alive
; and restarts them on crash/deploy. Install:
;
;   sudo cp deploy/supervisor/church-api.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;
; Adjust the paths (/var/www/church-api), `user`, and the `php` binary for your
; host. For horizontal scaling run several `reverb` nodes behind the TLS proxy
; (deploy/reverb/) with REVERB_SCALING_ENABLED=true so Redis pub/sub fans every
; broadcast out to all nodes.

[program:church-api-reverb]
process_name=%(program_name)s
command=php /var/www/church-api/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/church-api
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/church-api/storage/logs/reverb.log
stopwaitsecs=10

[program:church-api-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/church-api/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/church-api
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/church-api/storage/logs/queue.log
; Give in-flight jobs (e.g. tenant provisioning, media) time to finish on restart.
stopwaitsecs=3600

[program:church-api-scheduler]
process_name=%(program_name)s
command=php /var/www/church-api/artisan schedule:work
directory=/var/www/church-api
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/church-api/storage/logs/scheduler.log
stopwaitsecs=10

[group:church-api]
programs=church-api-reverb,church-api-queue,church-api-scheduler
